# pkg/scripts/postinstall
#!/bin/bash
##
## postinstall
##
## Executed automatically by macOS after the PKG is installed.
## Responsibilities:
##   1. Set correct permissions on scripts and LaunchDaemon plist
##   2. Inject the Intune client_secret into the System Keychain
##   3. Register and start the LaunchDaemon (which runs the migration)
##
## ── IMPORTANT ─────────────────────────────────────────────────────────────────
## The client_secret below must be rotated before each deployment.
## Never commit the real secret to the repos itory — use a placeholder here
## and inject the actual value through your secrets management pipeline
## (Jamf parameter, CI/CD variable, etc.) before building the PKG.
## ─────────────────────────────────────────────────────────────────────────────

# ── CONFIGURATION ─────────────────────────────────────────────────────────────
readonly COMPANY_NAME="ACME"                                          # <── change this
readonly BASE_DIR="/Library/Application Support/${COMPANY_NAME} MDM Migration"
readonly PLIST_NAME="com.acme.mdm.migration"                          # <── match your plist Label
readonly PLIST_PATH="/Library/LaunchDaemons/${PLIST_NAME}.plist"

# Intune App Registration — client_secret
# Replace with your actual secret (or inject via CI/CD pipeline)
INTUNE_SECRET="YOUR_CLIENT_SECRET_HERE"                               # <── change this

# Keychain identifiers — must match remover_intune.sh
readonly KEYCHAIN_SERVICE="MDMMigrationService"                       # <── change this
readonly KEYCHAIN_ACCOUNT="IntuneAuth"                                # <── change this
# ──────────────────────────────────────────────────────────────────────────────

echo "[postinstall] Starting post-installation..."

###############################################################################
# STEP 1 — Set permissions on LaunchDaemon plist
###############################################################################
echo "[postinstall] Setting permissions on LaunchDaemon plist..."

chmod 644 "${PLIST_PATH}"
chown root:wheel "${PLIST_PATH}"

echo "[postinstall] ✓ Plist permissions set"

###############################################################################
# STEP 2 — Set execute permissions on all scripts and jq binaries
###############################################################################
echo "[postinstall] Setting execute permissions on scripts..."

chmod +x "${BASE_DIR}/bin/"*.sh

if [[ -f "${BASE_DIR}/bin/jq-macos-arm64" ]]; then
    chmod +x "${BASE_DIR}/bin/jq-macos-arm64"
fi
if [[ -f "${BASE_DIR}/bin/jq-macos-amd64" ]]; then
    chmod +x "${BASE_DIR}/bin/jq-macos-amd64"
fi

echo "[postinstall] ✓ Permissions set"

###############################################################################
# STEP 3 — Inject client_secret into System Keychain
###############################################################################
echo "[postinstall] Injecting Intune credentials into System Keychain..."

# Preventive cleanup to avoid duplicate entries
/usr/bin/security delete-generic-password \
    -s "${KEYCHAIN_SERVICE}" \
    -a "${KEYCHAIN_ACCOUNT}" \
    /Library/Keychains/System.keychain 2>/dev/null || true

# Store the secret
/usr/bin/security add-generic-password \
    -a "${KEYCHAIN_ACCOUNT}" \
    -s "${KEYCHAIN_SERVICE}" \
    -w "${INTUNE_SECRET}" \
    -T "/usr/bin/security" \
    /Library/Keychains/System.keychain

if [[ $? -eq 0 ]]; then
    echo "[postinstall] ✓ Credentials stored in System Keychain"
else
    echo "[postinstall] ✗ Failed to store credentials" >&2
    exit 1
fi

###############################################################################
# STEP 4 — Register and start the LaunchDaemon
###############################################################################
echo "[postinstall] Registering LaunchDaemon..."

# Unload first in case a previous version is running
/bin/launchctl bootout system "${PLIST_PATH}" 2>/dev/null || true

# Bootstrap (register + start) — RunAtLoad=true will execute the migration
/bin/launchctl bootstrap system "${PLIST_PATH}"

if [[ $? -eq 0 ]]; then
    echo "[postinstall] ✓ LaunchDaemon registered and started"
else
    echo "[postinstall] ✗ Failed to register LaunchDaemon" >&2
    exit 1
fi

echo "[postinstall] ✓ Post-installation complete"
# macOS deletes this script automatically after installation
exit 0
